# examples/OECC/config/rc.local
#!/bin/sh -e
#

export DLC=/psc/dlc
export PATH=$DLC/bin:$PATH

cd /psc/wrk
proserve sports2020 -S 20000 -n 30

cd /psc/wrk/oepas1
/psc/wrk/oepas1/bin/tcman.sh start

exit 0
